# Source/NatServer/conf/pre_towebser.conf
start on stopped rc RUNLEVEL=[35]
stop on starting runlevel [01246]

respawn
exec /home/jjshe/ToWebServer
